# packages/lib/systemd/system/n3n-supernode.service
[Unit]
Description=n3n supernode process
Documentation=man:n3n-supernode(8)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=n3n
Group=n3n
RuntimeDirectory=n3n
ExecStart=/usr/sbin/n3n-supernode start
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=
