# packages/lib/systemd/system/n3n-edge.service
[Unit]
Description=n3n edge process
After=network-online.target nfw.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=
ExecStart=/usr/sbin/n3n-edge start
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=
