# packages/lib/systemd/system/supernode.service
[Unit]
Description=n3n supernode process
After=network-online.target syslog.target
Wants=network-online.target

[Service]
Type=simple
User=n3n
Group=n3n
ExecStart=/usr/sbin/supernode start -f
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=
